# .config/systemd/user/urxvtd.service
[Unit]
Description=Urxvt Daemon

[Service]
Type=simple
ExecStart=/usr/bin/urxvtd
Restart=always
RestartSec=10

[Install]
WantedBy=xorg.target
